# updateSystemSoftware/update-software.service
[Unit]
Description=Check for System SoftWare update
 
[Service]
Type=simple
ExecStart=/home/limo/.config/polybar/script/updateSoftWareAndSystem.sh
